# deps/res/linux/crontab
* * * * * /home/root/TLBB.Billing/bin/linux/scripts/restart_lb.sh

